# debian/kolibri-server.service
# This provides the Kolibri server service, started by invoking the old
# style /etc/init.d/kolibri-server service for backwards-compatibility.

[Unit]
Description=A high performance web server setup for Kolibri
After=local-fs.target remote-fs.target network-online.target
Wants=network-online.target

[Service]
Type=simple
#TimeoutStartSec=infinity
#TimeoutStopSec=infinity
ExecStart=/etc/init.d/kolibri-server start
ExecStop=/etc/init.d/kolibri-server stop

[Install]
WantedBy=multi-user.target
