# debian/kolibri-server.cron.weekly
#!/bin/sh

CONFIG_FILE=/etc/default/kolibri

[ -r $CONFIG_FILE ] || { echo "Configuration not found: $CONFIG_FILE" && exit 1 ;}
[ -r $CONFIG_FILE ] && . $CONFIG_FILE


[ "$KOLIBRI_USER"=="" ] || { echo "$KOLIBRI_USER not set" && exit 1 ;}

su $KOLIBRI_USER -c "cd /usr/share/kolibri-server/; python3 -c 'import kolibri_server_setup;kolibri_server_setup.delete_redis_cache()'"
